# demo/duplex/rp5-performance.service
[Unit]
Description=Set CPU Governor to Performance Mode
After=multi-user.target

[Service]
Type=oneshot
ExecStart=/bin/bash -c 'for i in 0 1 2 3; do echo performance > /sys/devices/system/cpu/cpu$i/cpufreq/scaling_governor; done'
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
